# web-app/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

s6-setuidgid www-data

php-fpm --nodaemonize