# web-app/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

#su www-data

s6-setuidgid www-data

exec /usr/local/sbin/php-fpm --nodaemonize