# web-app/services.d/nginx/run
#!/bin/sh

s6-setuidgid nginx

nginx -g 'daemon off;'