# web-app/services.d/nginx/run
#!/bin/sh

s6-setuidgid www-data

exec nginx -g 'daemon off;'